# salt/elife-bot/config/home-deployuser-elife-bot.cron
# every five minutes
*/5 * * * * cd /opt/elife-bot && /opt/elife-bot/scripts/run_cron.sh {{ env }}
